# supervisord-apache2.conf
[program:apache2]
command=/start.sh
numprocs=1
autostart=true
autorestart=true